# health_check/rippersnapper/nwsc2/rippersnapper.crontab
*/10 * * * * root test ! -f /etc/nolocal && /usr/sbin/rippersnapper enforce | logger -t rippersnapper
