# slurm/avail/avail_monitoring.crontab
*/5 * * * * lsfadmin flock -xn /var/tmp/get_slurm_stat_clusters /usr/sbin/get_slurm_stat_clusters.ncar
